# 99-hublink-usb.rules
# Enable detailed logging for all USB block devices
SUBSYSTEM=="block", KERNEL=="sd[a-z]*", ACTION=="add", RUN+="/usr/bin/logger -t hublink-udev 'USB Device Added - Kernel: %k, Name: %n, Number: %n, Label: $env{ID_FS_LABEL}, Type: $env{ID_FS_TYPE}'"

# Match USB storage partitions with HUBLINK label
SUBSYSTEM=="block", KERNEL=="sd[a-z]1", ACTION=="add", ENV{ID_FS_LABEL}=="HUBLINK", \
    MODE="0666", \
    OWNER="root", GROUP="root", \
    RUN+="/bin/sh -c 'if mountpoint -q /media/hublink-usb; then umount /media/hublink-usb; fi; \
                      rm -rf /media/hublink-usb; \
                      chmod 755 /media && \
                      /usr/local/bin/mount-hublink.sh %k'"

# Handle removal of HUBLINK devices
SUBSYSTEM=="block", KERNEL=="sd[a-z]1", ACTION=="remove", ENV{ID_FS_LABEL}=="HUBLINK", RUN+="/usr/local/bin/unmount-hublink.sh"

# Log removals
SUBSYSTEM=="block", KERNEL=="sd[a-z]*", ACTION=="remove", RUN+="/usr/bin/logger -t hublink-udev 'USB Device Removed - Kernel: %k'" 